# deploy/systemd/autops-llm-small.service
[Unit]
Description=AUTOPS Small LLM (llama.cpp, 短 prompt 智能路由目标，高并发)
# 大小模型分流：80% 简单问答走小模型(快、可高并发)，复杂推理才走大模型。
# 由后端 AIOps 智能路由按 prompt 长度选择端点（SMALL_LLM_PROMPT_THRESHOLD）。
After=network.target

[Service]
Type=simple
User=root
# 模型路径与端口按部署实际调整；后端通过 SMALL_LLM_URL 指向此端口
Environment=LLM_MODEL=/opt/models/Qwen3.5-0.8B-Q4_K_M.gguf
ExecStart=/usr/local/bin/llama-server -m ${LLM_MODEL} -c 2048 -n -1 -t 8 --host 0.0.0.0 --port 8003 --jinja --parallel 8

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
MemoryMax=8G

StandardOutput=journal
StandardError=journal
SyslogIdentifier=autops-llm-small

[Install]
WantedBy=multi-user.target
